# examples/crontab
# Standard 5-field cron format.
# Set TZ on the container to control local schedule interpretation.
# Non-git schedule: archive to a mounted data directory.
0 6 * * * /usr/local/bin/refbolt sync --all --verbose
